# dp-devel/stats/automodify.cron
0 0,3,6,9,12,15,18,21 * * * lynx -dump http://texts01.archive.org/~charlz/dproofreaders/tools/project_manager/automodify.php
