# dp-devel/stats/automodify.cron
0 0,3,6,9,12,15,18,21 * * * lynx -dump <<SITE_URL>>/tools/project_manager/automodify.php > <<CODE_DIR>>/stats/latest_releases.txt
