# dp-devel/stats/hourlystats.cron
01 * * * * /home/charlz/public_html/dproofreaders/stats/hourlystats.pl


