# dp-devel/stats/hourlystats.cron
01 * * * * <<SITE_DIR>>/stats/hourlystats.pl


