# dp-devel/stats/stats.cron
01 * * * * <<SITE_DIR>>/stats/hourlystats.pl
01 0 * * * <<SITE_DIR>>/stats/dailystats.pl
02 0 * * * <<SITE_DIR>>/stats/statsgraph.pl
0,15,30,45 * * * * lynx -dump <<SITE_URL>>/tools/project_manager/automodify.php

