# dp-devel/stats/tempstats.cron
*/5 * * * * /var/www/gutenberg/stats/statsgraph.pl
